# equestria-os-git-askpass/pkg/equestria-os-git-askpass/etc/profile.d/equestria-os-git-askpass.sh
export SSH_ASKPASS=/usr/bin/equestria-os-git-askpass
